# deploy/crosspenis.service
# /etc/systemd/system/crosspenis.service
# sudo systemctl daemon-reload && sudo systemctl enable crosspenis && sudo systemctl start crosspenis

[Unit]
Description=Crosspenis
After=network.target

[Service]
Type=simple
User=__RUN_USER__
WorkingDirectory=__APP_DIR__
Environment=PYTHONUNBUFFERED=1
ExecStart=__APP_DIR__/venv/bin/python __APP_DIR__/app.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
